# system/lib/minerstat.service
[Unit]
Description=minerstat-service
Requires=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
PrivateTmp=no
ExecStart=/bin/sh -c "screen -wipe; screen -A -m -d -S minerstat sh /config/minerstat/minerstat.sh"

[Install]
WantedBy=multi-user.target